# vivaldi/vivaldi-cookies.service
[Unit]
Description=Delete vivaldi cookies daily.

[Service]
ExecStart=/home/Rafa/config/vivaldi/delete-cookies.sh
Type=simple
